# op-return-analytics.service
[Unit]
Description=OP_RETURN Analytics Dashboard
After=network.target

[Service]
Type=simple
User=ben
WorkingDirectory=/home/ben/projects/op-return-analytics
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=5

Environment=PORT=8083
#Environment=OPRETURN_DB=/path/to/invoices.sqlite

[Install]
WantedBy=multi-user.target
